# qdrant.service
[Unit]
Description=Qdrant Vector Database for Timus
After=network.target local-fs.target
Wants=network.target local-fs.target

[Service]
Type=simple
User=fatih-ubuntu
WorkingDirectory=/home/fatih-ubuntu/dev/timus
ExecStart=/home/fatih-ubuntu/dev/timus/scripts/start_qdrant_server.sh
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=qdrant
TimeoutStartSec=45
TimeoutStopSec=20

[Install]
WantedBy=multi-user.target
